# systemd/oktoberfest-bot.service.example
[Unit]
Description=Oktoberfest Tent Reservation Monitor
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/oktoberfest-bot
ExecStart=/usr/bin/python3 -m oktoberfest_bot.main
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
